# dotfiles/dot_config/autostart/1password.desktop
[Desktop Entry]
Type=Application
Name=1Password
Exec=1password --silent
Terminal=false
X-GNOME-Autostart-enabled=true
NoDisplay=true
